# pi400kb.service
[Unit]
Description=pi400kb USB OTG Keyboard & Mouse forwarding

[Service]
ExecStart=/usr/bin/pi400kb
User=root   
Group=root   
Type=simple
Restart=on-failure

[Install]
WantedBy=multi-user.target
